# assets/conf/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=nginx -c /usr/local/nginx/conf/nginx.conf
